# files/default/supervisord.d/cloudinit-passwd.ini
[program:cloudinit-passwd]
command=/usr/bin/python /opt/cloudinit-metaserv/cloud-passwd.py
directory=/opt/cloudinit-metaserv
user=root
autostart=true
autorestart=true
redirect_stderr=true
