# files/default/supervisord.d/cloudinit-meta.ini
[program:cloudinit-meta]
command=/usr/bin/python /opt/cloudinit-metaserv/cloud-meta.py
directory=/opt/cloudinit-metaserv
user=root
autostart=true
autorestart=true
redirect_stderr=true
